# cron.d/backup-crontab
0 */4 * * * /opt/minecraft/server/backup.sh
